# systemd/pa-email-retention.timer
[Unit]
Description=PA_V2 Email Retention Cleanup Timer
Requires=pa-email-retention.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
