# systemd/pa-email-check-hourly.timer
[Unit]
Description=Run email check every hour (alternative to minutely)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=pa-email-check.service

[Install]
WantedBy=timers.target
