# systemd/pa-email-check.service
[Unit]
Description=PA: one-shot email check (Gmail + Outlook)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mentorius
WorkingDirectory=/home/mentorius/AI_Services/PA_V2
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/pa_v2/secrets.env
# Use the virtualenv python for stability; installer will replace if needed
ExecStart=/home/mentorius/AI_Services/PA_V2/.venv/bin/python -m src.jobs.email_check
